# crontab
0   12  *   *   6  notifier.sh access.log error.log 10 10
